# infra/docker/s6-rc.d/web/run
#!/usr/bin/with-contenv bash
# uvicorn — FastAPI + bundled SvelteKit
cd /srv/apps/web
exec s6-setuidgid "${PUID:-1000}:${PGID:-1000}" \
  uvicorn setvault_web.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --proxy-headers \
    --forwarded-allow-ips '*'
