# infra/docker/s6-rc.d/init-db/up
/usr/bin/with-contenv bash /etc/setvault/init-db.sh
