# infra/docker/s6-rc.d/worker/run
#!/usr/bin/with-contenv bash
# RQ worker — transcode / waveform / enrichment / bulk jobs
cd /srv
exec s6-setuidgid "${PUID:-1000}:${PGID:-1000}" \
  python -m setvault_worker.entrypoint
